# share/debian/dogbot.cron
*/1 * * * * /etc/init.d/dogbot reload && /etc/init.d/dogbot restart
