# cron/crontab
*/1 * * * * /etc/init.d/feedeobot reload && /etc/init.d/feedeobot restart
