# cron/crontab
*/1 * * * * /root/feedeobot/cron/reload
