# cron/crontab
*/5 * * * * /root/feedeobot/cron/reload 2>&1 >/dev/null
